# etc/xinetd.d/atmail
service atmail
{
        socket_type     = stream
        protocol        = tcp
        user            = atmail
        wait            = no
        server          = /home/atmail/atmail
}
