# etc/xinetd.d/babyecho
service babyecho
{
        socket_type     = stream
        protocol        = tcp
        user            = babyecho
        wait            = no
        server          = /home/babyecho/babyecho
}
