# etc/xinetd.d/bookworm
service bookworm
{
        socket_type     = stream
        protocol        = tcp
        user            = bookworm
        wait            = no
        server          = /home/bookworm/bookworm
}
