# etc/xinetd.d/replayme
service replayme
{
        socket_type     = stream
        protocol        = tcp
        user            = replayme
        wait            = no
        server          = /home/replayme/replayme
}

