# etc/xinetd.d/trouver
service trouver
{
        socket_type     = stream
        protocol        = tcp
        user            = trouver
        wait            = no
        server          = /home/trouver/trouver
}
