# etc/xinetd.d/wdub
service wdub
{
        socket_type     = stream
        protocol        = tcp
        user            = wdub
        wait            = no
        server          = /home/wdub/wdub
}

